# configs/trade/paper_trade_daily.crontab.example
SHELL=/bin/bash
PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
CRON_TZ=Asia/Taipei

# Preview run before market open.
# This refreshes today's preview order list using the latest completed trade date.
30 08 * * 1-5 /home/nas2/Personal/Hank/qlib-public/scripts/trade/run_paper_trade_daily.sh

# Close run after market close.
# This refreshes fills, paper_state, nav_history, and order_fill_comparison.csv.
10 16 * * 1-5 /home/nas2/Personal/Hank/qlib-public/scripts/trade/run_paper_trade_daily.sh

# Example override:
# 10 16 * * 1-5 PAPER_TRADE_CONFIG=/home/nas2/Personal/Hank/qlib-public/configs/trade/paper_trading.alpha158_lgb_run11_tplus.example.json /home/nas2/Personal/Hank/qlib-public/scripts/trade/run_paper_trade_daily.sh
